# backend/deploy/reporting-api-gunicorn.service
[Unit]
Description=Reporting API (FastAPI via Gunicorn+Uvicorn)
After=network.target

[Service]
Type=simple
# Change to the absolute path of your backend directory
WorkingDirectory=/users/mohammed/Documents/Bayan/backend/
# Ensure run_prod_gunicorn.sh is executable
ExecStart=/users/mohammed/Documents/Bayan/backend/run_prod_gunicorn.sh
# Optionally override environment here or in an EnvironmentFile
Environment=RUN_SCHEDULER=0
Environment=HOST=127.0.0.1
Environment=PORT=8000
# Optionally control worker count
# Environment=WEB_CONCURRENCY=6
# Run as a non-root service account
User=www-data
Group=www-data
# Restart policy
Restart=always
RestartSec=5
# Increase file descriptor limit for high concurrency
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
